# cdh/supervisor_edgenode.ini
[program:jupyterhub]
environment=PATH=/usr/local/anaconda/python3/bin:/usr/local/sbin:/sbin:/bin:/usr/sbin:/usr/bin:/root/bin
command=jupyterhub --no-ssl --debug --log-level=0 --log-file=/var/log/jupyterhub/jupyter.log --db=/etc/jupyterhub/jupyterhub.sqlite --config=/etc/jupyterhub/jupyterhub_config.py
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/jupyterhub/jupyterhub-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stdout_events_enabled=true
stderr_logfile=/var/log/jupyterhub/jupyterhub-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
stderr_events_enabled=true

[program:linuxusersync]
command=/etc/supervisord.d/scripts/syncpass.sh
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/linuxusersync/linuxusersync-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stdout_events_enabled=true
stderr_logfile=/var/log/linuxusersync/linuxusersync-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
stderr_events_enabled=true

[program:webconsole]
command=/bin/forever start /opt/wetty/app.js -p 42000
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/webconsole/webconsole-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stdout_events_enabled=true
stderr_logfile=/var/log/webconsole/webconsole-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
stderr_events_enabled=true

[program:httpd]
command=/sbin/service httpd start
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/httpd/httpd-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stdout_events_enabled=true
stderr_logfile=/var/log/httpd/httpd-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
stderr_events_enabled=true

[program:ftpd]
command=/sbin/service vsftpd start
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/ftpd/ftpd-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stdout_events_enabled=true
stderr_logfile=/var/log/ftpd/ftpd-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
stderr_events_enabled=true

[program:rstudio]
command=/sbin/rstudio-server start
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/rstudio/rstudio-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stdout_events_enabled=true
stderr_logfile=/var/log/rstudio/rstudio-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
stderr_events_enabled=true
